# openclaw-skill/twinbox-openclaw-bridge.timer
[Unit]
Description=Twinbox OpenClaw cron bridge poller timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=twinbox-openclaw-bridge.service

[Install]
WantedBy=timers.target
