# integrations/openclaw/twinbox-openclaw-bridge.service
[Unit]
Description=Twinbox OpenClaw cron bridge poller
After=network-online.target openclaw-gateway.service
Wants=network-online.target openclaw-gateway.service

[Service]
Type=oneshot
EnvironmentFile=%h/.twinbox/twinbox-openclaw-bridge.env
WorkingDirectory=%h
ExecStart=/bin/bash -lc 'cd "${TWINBOX_STATE_ROOT:-${TWINBOX_CANONICAL_ROOT}}" && "${TWINBOX_CODE_ROOT}/scripts/twinbox_openclaw_bridge_poll.sh" --format json --openclaw-bin "${OPENCLAW_BIN:-openclaw}"'
